# AtraLInux/_Scripts/mnt-mail.automount
[Unit]
Description=Automount Fresh series

[Automount]
Where=/mnt/mail
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target
