# AtraLInux/_Scripts/mnt-consplusregion.automount
[Unit]
Description=Automount Fresh series

[Automount]
Where=/mnt/consplusregion
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target
