# Scripts/mnt-base.automount
[Unit]
Description=Automount Fresh series

[Automount]
Where=/mnt/base
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target
